# scraper/jobs/crontab
30 * * * * ~/env/bin/python ~/current/vgn.py
